# hamwing-gpio-init.service
[Unit]
Description=HamWing GPIO safe-state initializer
# Run as early as possible, before any radio scripts.
After=sysinit.target
Before=multi-user.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/pi-user/pi_gpio_init.py
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
